# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid

# nginx needs to start as root so it can bind :80; its worker processes
# drop to 'www-data' per the stock nginx.conf user directive.
[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/var/log/nginx.stdout.log
stderr_logfile=/var/log/nginx.stderr.log

# FastAPI runs as the unprivileged 'appuser'. All user-uploaded scripts
# execute under this uid, sandboxed further by firejail (see
# backend/services/executor.py).
[program:fastapi]
command=uvicorn backend.main:app --host 127.0.0.1 --port 8000
directory=/app
user=appuser
autostart=true
autorestart=true
stdout_logfile=/var/log/fastapi.stdout.log
stderr_logfile=/var/log/fastapi.stderr.log

# Next.js SSR runs as the unprivileged 'appuser'.
[program:nextjs]
command=node /app/frontend/node_modules/.bin/next start --port 3000
directory=/app/frontend
user=appuser
autostart=true
autorestart=true
stdout_logfile=/var/log/nextjs.stdout.log
stderr_logfile=/var/log/nextjs.stderr.log
