# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid

[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/var/log/nginx.stdout.log
stderr_logfile=/var/log/nginx.stderr.log

[program:fastapi]
command=uvicorn backend.main:app --host 127.0.0.1 --port 8000
directory=/app
autostart=true
autorestart=true
stdout_logfile=/var/log/fastapi.stdout.log
stderr_logfile=/var/log/fastapi.stderr.log

[program:nextjs]
command=node /app/frontend/node_modules/.bin/next start --port 3000
directory=/app/frontend
autostart=true
autorestart=true
stdout_logfile=/var/log/nextjs.stdout.log
stderr_logfile=/var/log/nextjs.stderr.log
